# external/board/rpi3/rootfs_overlay/etc/init.d/S90sgw
#!/bin/sh
#
# chkconfig: 35 90 12
# description: sgw service
#

# Start the service sgw
start() {
        initlog -c "echo -n Starting sgw service: "
        /usr/bin/sgw &
        ### Create the lock file ###
        touch /var/lock/subsys/sgw
        success $"sgw server startup"
        echo
}

# Restart the service sgw
stop() {
        initlog -c "echo -n Stopping sgw service: "
        killproc sgw
        ### Now, delete the lock file ###
        rm -f /var/lock/subsys/sgw
        echo
}

### main logic ###
case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  status)
        status sgw
        ;;
  restart|reload|condrestart)
        stop
        start
        ;;
  *)
        echo $"Usage: $0 {start|stop|restart|reload|status}"
        exit 1
esac

exit 0
